# httpd/lib/systemd/system/httpd.service
[Unit]
Description=CERN httpd
[Service]
Type=Simple
RestartSec=1
ExecStart=/usr/bin/httpd
[Install]
WantedBy=multi-user.target
